# roles/chia_rpc/templates/chia_rpc.service.j2
[Unit]
Description=chia_rpc unit service

[Service]
User=root
ExecStart="python3 {{ chia_rpc_clone_path}}/main.py"
Restart=always

[Install]
WantedBy=multi-user.target
